# glomp.service
[Unit]
Description=reload my webpage when i push a new version
After=multi-user.target

[Service]
ExecStart=
